# supervisor.d/copyonstart.conf
[program:copyonstart]
command = /copyscript.sh
startsecs = 0
autorestart = false
startretries = 1
